# steamdeck_brain/systemd/scholomance-brain.service
[Unit]
# ─────────────────────────────────────────────────────────────────────────────
# scholomance-brain.service — SteamDeck Brain HTTP daemon (no-llm by default)
# ─────────────────────────────────────────────────────────────────────────────
# Runs brain_daemon.py (http://127.0.0.1:9090) as a one-shot user service.
# DivTube connects to this daemon over HTTP; it does not own its lifecycle.
#
# No-llm mode: the daemon serves the full ForceField + SCDNA pipeline without
# loading or calling Ollama. No 7B model is touched, no extra RAM is pinned.
# If you want the Ollama-backed LLM bridge instead, run the daemon manually:
#   ./launch-brain-stack.sh   (or set SCHOLOMANCE_WITH_OLLAMA=1)
#
# Restart is disabled and auto-start is opt-in. Manual launches control the
# lifecycle; the daemon will not fight itself for the port.
#
# This file is a TEMPLATE: @PYTHON@, @SCRIPT_DIR@ and @BRAIN_DAEMON@ are
# substituted by install.sh.
Description=Scholomance SteamDeck Brain daemon (no-llm: ForceField + SCDNA)
After=network-online.target
# Must live in [Unit] — systemd ignores it under [Service].
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=@SCRIPT_DIR@
# Optional overrides: drop KEY=VALUE lines (BRAIN_MODEL, BRAIN_DB, BRAIN_PORT,
# BRAIN_PERSONALITY, BRAIN_EXTRA_ARGS) into this file. "-" = ok if missing.
EnvironmentFile=-%h/.config/scholomance-brain.env
Environment=BRAIN_MODEL=qwen2.5-coder:7b
Environment=BRAIN_PORT=9090
Environment=BRAIN_DB=%h/.substrate/memory.sqlite
ExecStart=@PYTHON@ @BRAIN_DAEMON@ --no-llm --model ${BRAIN_MODEL} --port ${BRAIN_PORT} --db ${BRAIN_DB} $BRAIN_EXTRA_ARGS
# No auto-restart. Manual launches own the lifecycle; the daemon will not
# fight itself for the port and will not pin RAM for a model.
Restart=no

[Install]
# No default.target — the service does NOT auto-start at boot/login.
# Enable it explicitly with:  systemctl --user enable scholomance-brain.service
